# .config/systemd/user/blueman.service
[Unit]
Description=Blueman bluetooth manager

[Service]
ExecStart=/usr/bin/blueman-applet

[Install]
WantedBy=i3.target
