# .config/systemd/user/i3.target
[Unit]
Description=i3 background processes
