# .config/systemd/user/unclutter.service
[Unit]
Description=Unclutter

[Service]
Environment=DISPLAY=:0
ExecStart=/usr/bin/unclutter

[Install]
WantedBy=i3.target
